# Microproyecto1/puppetNodo2/modules/nodeconfig/files/node.service
[Service]
ExecStart=/usr/bin/node /home/vagrant/app/index.js 3004
Restart=on-failure
WorkingDirectory=/home/vagrant/app