# Microproyecto1/puppetNodo1/modules/baseconfig/files/node.service
[Service]
ExecStart=/usr/bin/node /home/vagrant/app/index.js 3002
Restart=on-failure