# Microproyecto1/puppetNodo1/modules/consulconfig/files/consul.service
[Service]
ExecStart=/usr/bin/consul agent -config-file=/usr/local/etc/consul/config.json
Restart=on-failure